# okweb/deploy/web.hujiangtao.cn.ini
[program:web.hujiangtao.cn]
environment=PYTHONPATH='/home/jiangtao/django-web/okweb:/home/jiangtao/django-web/.venv/lib/python3.6/site-packages/'
command=pipenv run python Server.py
directory=/home/jiangtao/django-web/okweb
user=root
stopwaitsecs=60
stopasgroup=true
killasgroup=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/jiangtao/web.hujiangtao.cn.log
loglevel=info